# auditor/systemd/tine-plugin-auditor.timer
[Unit]
Description=Poll Tine plugin submissions every five minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=5m
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
